# docs/supervisor/example-psbot-worker.conf
[program:psbot-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/psbot/artisan queue:work worker --queue=bannerWorker,afkWorker,channelRemoverWorker,policeWorker --once --sleep=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/psbot-worker.log
stopwaitsecs=3600
startsecs=0
